# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server

[program:yaas]
command=gulp && pm2 start lib/index.js --no-daemon
