# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server

[program:yaas]
command=run.sh