# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server